# jarvis-secretsd/jarvis-secretsd/jarvis-secretsd.service
[Unit]
Description=Jarvis Secrets Daemon
Documentation=https://github.com/jarvis/jarvis-secretsd
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/jarvis-secretsd --config /etc/jarvis-secretsd/config.toml

# User and group
User=jarvis
Group=jarvis

# Security hardening
AmbientCapabilities=
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Read/write paths
ReadWritePaths=/opt/jarvis

# Capability bounding
CapabilityBoundingSet=

# Restart policy
Restart=on-failure
RestartSec=5s

# Environment
Environment="RUST_LOG=info"

[Install]
WantedBy=multi-user.target
